# templates/etcd.upstart.j2
description "etcd service registry"

start on started
stop on shutdown

exec /usr/bin/etcd -d /var/data/etcd/config -n config -c {{ ip }}:4001
